# systemd/ransomeye-linux-agent.service
# Path and File Name: /home/ransomeye/rebuild/systemd/ransomeye-linux-agent.service
# Author: nXxBku0CKFAJCBN3X1g3bQk7OxYQylg8CMw1iGsq7gU
# Details: Systemd service unit for RansomEye Linux Agent with runtime hardening

[Unit]
Description=RansomEye Linux Agent - Host telemetry collection
After=network.target
Requires=network.target
ConditionPathExists=/opt/ransomeye/linux_agent/bin/ransomeye_linux_agent
ConditionPathExists=/opt/ransomeye/linux_agent/.install_receipt.json

[Service]
Type=notify
# Rootless runtime enforcement - run as ransomeye user
User=ransomeye
Group=ransomeye
WorkingDirectory=/opt/ransomeye/linux_agent
RuntimeDirectory=ransomeye/linux_agent
StateDirectory=ransomeye/linux_agent
ExecStart=/opt/ransomeye/linux_agent/bin/ransomeye_linux_agent
# Watchdog timer - service must send heartbeat every 30 seconds
WatchdogSec=30
# Restart on failure (anti-kill protection)
Restart=always
RestartSec=10
# Restart on watchdog timeout
RestartForceExitStatus=143
# Maximum restart attempts before escalation
StartLimitIntervalSec=300
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

# Resource limits (prevent swap usage)
LimitNOFILE=65536
LimitNPROC=4096
# NO swap allowed - hard limit
MemoryMax=2G
MemorySwapMax=0

# Security hardening - Rootless runtime enforcement
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ransomeye/linux_agent /var/lib/ransomeye/linux_agent /run/ransomeye/linux_agent

# Additional hardening
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateDevices=false
# Allow network access but restrict listeners
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Capability-based privileges (no root required)
CapabilityBoundingSet=CAP_SYS_PTRACE CAP_DAC_READ_SEARCH CAP_NET_RAW CAP_NET_ADMIN
AmbientCapabilities=
PrivateUsers=false

# Memory policy (NO swap)
# Systemd will enforce MemorySwapMax=0 above
# Additional protection via mlock (if needed by binary)

# Environment variables
EnvironmentFile=-/opt/ransomeye/linux_agent/config/environment.conf

# Kill settings - fail-closed on kill attempts
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
# Send alert on unexpected termination
TimeoutStartSec=60

# Privilege downgrade is handled by the binary itself
# Binary starts as root (SetUID) and drops to unprivileged user

[Install]
WantedBy=multi-user.target

